# dev/FreeBSD/feature/programming-languages/java/sonarqube.feature/setup/files/etc/rc.d/sonarqube
#!/sbin/openrc-run
PIDFILE=/run/sonarqube.pid
_SONARQUBE_VERSION=7.1
_SONARQUBE_EXECUTABLE=/opt/sonarqube.com/sonarqube-$_SONARQUBE_VERSION/bin/linux-x86-64/sonar.sh
start() {
	ebegin "Starting sonarqube"
	#start-stop-daemon --start --quiet -u sonarqube -g sonarqube-b -m -p ${PIDFILE} --exec sonarqube
	sudo -u sonarqube $_SONARQUBE_EXECUTABLE start
	eend $? 
}
stop() {
	ebegin "Stopping sonarqube"
	#start-stop-daemon --stop --quiet --name sonarqube
	$_SONARQUBE_EXECUTABLE stop
	eend $?
}
status() {
	local _retval
	$_SONARQUBE_EXECUTABLE status>/dev/null
	_retval=$?
	if [ ${_retval} = '0' ]; then
		einfo 'status: started'
		return 0
	else
		einfo 'status: stopped'
		return 3
	fi
}
